# res/e-mission-server/scripts/crontab
@daily root /usr/src/app/run.cronjob.sh
