# e-mission-server/crontab
@daily root /usr/src/app/dk_run_cronjob.sh
